# conf/nodejs.conf
[supervisord]
nodaemon=true

[program:nodesprocess]
command=/usr/bin/nodejs /app/app.js
directory=/app
autostart=true
autorestart=true
startretries=3