# conf/nodejs.conf
[supervisord]
nodaemon=true

[program:nodesprocess]
command=/bin/su -c "/usr/bin/nodejs /app/app.js" - bakdyr
directory=/app
autostart=true
autorestart=true
startretries=3